# deprecated/docker/allinone/src/goreplay.ini
[program:goreplay]
command=/opt/goreplay/bin/gor --input-raw %(ENV_VPC_MIRROR_DEVICE)s:%(ENV_APP_PORT)s --input-raw-track-response --input-raw-bpf-filter "(dst port "%(ENV_APP_PORT)s") or (src port "%(ENV_APP_PORT)s")" --output-resurface "%(ENV_USAGE_LOGGERS_URL)s" --output-resurface-rules "%(ENV_USAGE_LOGGERS_RULES)s"
autostart=false
