# src/goreplay.ini
[program:goreplay]
command=/opt/goreplay/bin/gor --input-raw :%(ENV_APP_PORT)s --input-raw-track-response --output-resurface "%(ENV_USAGE_LOGGERS_URL)s" --output-resurface-rules "%(ENV_USAGE_LOGGERS_RULES)s"
autostart=false
